# dka-postgresql/postgres17/config/cron.d/backup
#!/bin/bash

# Set database credentials
DB_USER="postgres"
# Set default values if environment variables are not provided
DB_PASSWORD=${DKA_DB_PASSWORD:-postgres}
BACKUP_DIR="/backup"
# shellcheck disable=SC1001
DATE=$(date +\%d-\%m-\%Y_\%H-\%M-\%S)
BACKUP_FILE="$BACKUP_DIR/schedule_all_databases_$DATE.sql.gz"

# Optional S3/Wasabi upload settings
# Set DKA_S3_UPLOAD_ENABLE=true to enable automatic upload after backup
DKA_S3_UPLOAD_ENABLE=${DKA_S3_UPLOAD_ENABLE:-false}
DKA_S3_BUCKET=${DKA_S3_BUCKET:-}
DKA_S3_PATH=${DKA_S3_PATH:-}
DKA_S3_ENDPOINT=${DKA_S3_ENDPOINT:-https://s3.wasabisys.com}
DKA_S3_REGION=${DKA_S3_REGION:-us-east-1}

# Create a backup directory if it doesn't exist
mkdir -p $BACKUP_DIR

# Set the PostgreSQL password for the session
export PGPASSWORD=$DB_PASSWORD

# If DKA_DB_NAME is set, backup single database using pg_dump; otherwise fallback to pg_dumpall
if [ -n "${DKA_DB_NAME:-}" ]; then
    DB_NAME_SINGLE=${DKA_DB_NAME}
    BACKUP_FILE="$BACKUP_DIR/${DB_NAME_SINGLE}_$DATE.sql.gz"
    echo "Backing up single database: $DB_NAME_SINGLE"
    pg_dump -U "$DB_USER" -d "$DB_NAME_SINGLE" \
        --no-owner \
        --no-privileges \
        --clean \
        --if-exists \
        --format=p | gzip -9 > "$BACKUP_FILE"
else
    echo "Backing up all databases..."

    # Create a backup of all databases, including users, ownership, and roles, and compress it into a .sql.gz file using pg_dumpall
    pg_dumpall -U $DB_USER --clean --no-owner | gzip -9 > "$BACKUP_FILE"
fi

if [ $? -eq 0 ]; then
    echo "Backup completed successfully: $BACKUP_FILE"
    # If S3 upload is enabled, try to upload using AWS CLI (supports Wasabi via endpoint)
    if [ "$DKA_S3_UPLOAD_ENABLE" = "true" ]; then
        if [ -z "$DKA_S3_BUCKET" ]; then
            echo "S3 upload enabled but DKA_S3_BUCKET is not set. Skipping upload."
        else
            # ensure path formatting
            UPLOAD_TARGET="s3://$DKA_S3_BUCKET"
            if [ -n "$DKA_S3_PATH" ]; then
                UPLOAD_TARGET="$UPLOAD_TARGET/${DKA_S3_PATH#/}"
            fi

            if command -v aws >/dev/null 2>&1; then
                echo "Uploading $BACKUP_FILE -> $UPLOAD_TARGET (endpoint=$DKA_S3_ENDPOINT)"
                aws s3 cp "$BACKUP_FILE" "$UPLOAD_TARGET/" --endpoint-url "$DKA_S3_ENDPOINT" --region "$DKA_S3_REGION"
                if [ $? -eq 0 ]; then
                    echo "Upload successful."
                else
                    echo "Upload failed."
                fi
            else
                echo "aws CLI not found in container. Install AWS CLI to enable S3/Wasabi uploads."
            fi
        fi
    fi
else
    echo "Backup failed."
fi
